# recipes-graphics/dmlauncher/dmlauncher/dmlauncher.service
[Unit]
Description=DM Launcher

[Service]
Type=simple
ExecStart=/usr/bin/dmlauncher
Restart=always

[Install]
WantedBy=multi-user.target
